# Solution Items/GhandiBot.conf
[program:GhandiBot]
command = #(directory)/GhandiBot
autostart = true