# Solution Items/GhandiBotStaging.conf
[program:GhandiBotStaging]
command = #(directory)/GhandiBot
autostart = true