# Solution Items/GhandiBotStaging.conf
[program:GhandiBotStaging]
command = #(directory)/GhandiBotStaging
autostart = true