# Solution Items/GhandiBotStaging.conf
[program:GhandiBotStaging]
command = #(directory)
autostart = true